# server/pianobar.service
[Service]
Type=forking
PIDFile=/var/run/pianobar.pid
ExecStart=/usr/sbin/pianobar
